# .config/systemd/user/battery_notify.timer
[Unit]
Description=run battery notifier every 5min

[Timer]
OnBootSec=30s
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=default.target

